# deploy/babgwe.service
[Unit]
Description=Babgwe daily Slack lunch recommendation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=babgwe
Group=babgwe
WorkingDirectory=/opt/babgwe
EnvironmentFile=/etc/babgwe/babgwe.env
ExecStart=/opt/babgwe/.venv/bin/python -m babgwe

