# deploy/bapratustra.timer
[Unit]
Description=Run Bapratustra at 11:00 KST on weekdays

[Timer]
OnCalendar=Mon..Fri *-*-* 11:00:00 Asia/Seoul
Unit=bapratustra.service

[Install]
WantedBy=timers.target
