# deploy/bapratustra-failure@.service
[Unit]
Description=Notify Slack that %i failed
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=bapratustra
Group=bapratustra
WorkingDirectory=/opt/bapratustra
EnvironmentFile=/etc/bapratustra/bapratustra.env
ExecStart=/opt/bapratustra/.venv/bin/python -m bapratustra --notify-systemd-failure %i
TimeoutStartSec=60
