# deploy/bapratustra-slack.service
[Unit]
Description=Bapratustra Slack interaction service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bapratustra
Group=bapratustra
WorkingDirectory=/opt/bapratustra
EnvironmentFile=/etc/bapratustra/bapratustra.env
ExecStart=/opt/bapratustra/.venv/bin/python -m bapratustra --run-slack-service
Restart=on-failure
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
