# deploy/bapratustra-alpha-backup.service
[Unit]
Description=Back up the Bapratustra alpha database
After=bapratustra-alpha.service
OnFailure=bapratustra-failure@%n.service

[Service]
Type=oneshot
User=bapratustra
Group=bapratustra
WorkingDirectory=/opt/bapratustra
EnvironmentFile=/etc/bapratustra/bapratustra.env
Environment=PYTHON_DOTENV_DISABLED=1
ExecStart=/opt/bapratustra/.venv/bin/python -m bapratustra --backup-alpha
TimeoutStartSec=120
